# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# Enable `supervisorctl` to talk to this supervisord over a local unix
# socket — lets `docker exec <container> supervisorctl restart stra2us`
# bounce the app without restarting the whole container (which would
# also bounce redis and lose its in-process pool/cache state).
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface.make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:redis]
command=redis-server --dir /data --appendonly yes
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.out.log
stderr_logfile=/var/log/supervisor/redis.err.log

[program:stra2us]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 8153
directory=/app/src
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/stra2us.out.log
stderr_logfile=/var/log/supervisor/stra2us.err.log
environment=PYTHONPATH="/app/src"
