# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# Enable `supervisorctl` to talk to this supervisord over a local unix
# socket — lets `docker exec <container> supervisorctl restart stra2us`
# bounce the app without restarting the whole container (which would
# also bounce redis and lose its in-process pool/cache state).
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
# Note: colon (not dot) between module and attribute is required by
# current setuptools/pkg_resources EntryPoint.parse. The older dotted
# form `supervisor.rpcinterface.make_main_rpcinterface` is rejected by
# python:3.11-slim's bundled pkg_resources with "ModuleNotFoundError:
# 'supervisor.rpcinterface' is not a package".
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:redis]
command=redis-server --dir /data --appendonly yes
autostart=true
autorestart=true
# Redirect to docker stdout/stderr so `docker logs` and `tools/stage
# logs` see redis output without `docker exec` + cat'ing files inside
# the container. supervisord doesn't rotate /dev/{stdout,stderr}, so
# *_logfile_maxbytes=0 disables its own rotation logic.
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:stra2us]
# Wrapped in `sh -c '... exec ...'` so we can default-expand
# $STRA2US_WORKERS (supervisord's own %(ENV_X)s interpolation
# has no `:-default` syntax). The `exec` is load-bearing —
# without it supervisord supervises the shell, not uvicorn, and
# signals don't propagate cleanly.
#
# --workers: one process per worker, each running its own
#   asyncio event loop on its own core. Default 4 — adjust via
#   STRA2US_WORKERS env var in compose if the host's CPU count
#   warrants it. All workers share the same baked
#   /etc/stra2us/admin_session_secret so admin cookies
#   validate uniformly across them.
# --limit-concurrency: per-worker in-flight request cap.
#   Backpressure floor — under a load spike, workers return
#   503s before queueing unbounded and OOMing. 1000 per worker
#   × 4 workers = 4000 in-flight requests before any rejection,
#   well above expected steady state.
command=sh -c 'exec python3 -m uvicorn main:app --host 0.0.0.0 --port 8153 --workers ${STRA2US_WORKERS:-4} --limit-concurrency ${STRA2US_LIMIT_CONCURRENCY:-1000}'
directory=/app/src
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/app/src"
